# wschat.conf
description "wschat"

start on runlevel [2534]
stop on runlevel [!2534]

script
    cd /home/wagnerpatriota/www/chat
    node server.js
end script

respawn
